# ops/hardening/srv-backend-watchdog.timer
[Unit]
Description=Run SRV production backend liveness watchdog every minute

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
AccuracySec=5s
Persistent=true
Unit=srv-backend-watchdog.service

[Install]
WantedBy=timers.target
